# roles/motion/templates/supervisor_motion.conf.j2
[program:motion]
command=motion -m -c /home/pi/.motion/motion.conf
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/motion.log
user=pi
